# node-red/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: Node-RED
# Runs the Nginx daemon
# ==============================================================================

# Wait for Node-RED to become available
bashio::net.wait_for 46836 localhost 300

bashio::log.info "Starting NGinx..."
exec nginx
